# deploy/rw-ast-tools.service
[Unit]
Description=rw-ast-tools MCP Server (Persistent Daemon)
Documentation=https://github.com/stephanos8926-lgtm/ast-tools
After=network.target

[Service]
Type=simple
ExecStart=/home/sysop/Workspaces/ast-tools/.venv/bin/ast-tools-server --mode daemon --foreground
Restart=on-failure
RestartSec=5
WatchdogSec=120
StartLimitBurst=5
StartLimitIntervalSec=300
MemoryMax=1G
Environment=AST_TOOLS_MODE=daemon
Environment=AST_TOOLS_LOG_LEVEL=info

# Security hardening
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/home/sysop/.cache/rw-ast-tools /home/sysop/.config/rw-ast-tools /home/sysop/Workspaces /tmp

[Install]
WantedBy=default.target